# ai_news_scraper/deploy/ai-news-pipeline.service
[Unit]
Description=AI News Pipeline — Multi-Interest Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ai-news-pipeline
Group=ai-news-pipeline
WorkingDirectory=/opt/ai-news-pipeline
EnvironmentFile=/opt/ai-news-pipeline/.env
ExecStart=/opt/ai-news-pipeline/venv/bin/python -m src.main \
  --config /opt/ai-news-pipeline/config/ \
  --serve \
  --log-file /opt/ai-news-pipeline/logs/pipeline.log
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
